# misc/docker/service/services.d/rrr/finish
#!/bin/sh

echo "Service rrr exited with code $1, restarting"

exit 0
