# misc/systemd/rrr.service
[Unit]
Description=Read Route Record 0.18

[Service]
ExecStart=/usr/bin/rrr /etc/rrr.conf.d/
KillSignal=SIGINT
User=rrr
Group=rrr

[Install]
WantedBy=multi-user.target
