# misc/systemd/rrr.service
[Unit]
Description=Read Route Record 1.18-0

[Service]
ExecStart=/usr/local/bin/rrr -s -e /usr/local/etc/rrr.conf /usr/local/etc/rrr.conf.d/
KillSignal=SIGTERM
User=daemon
Group=daemon

[Install]
WantedBy=multi-user.target
