# misc/systemd/rrr.service
[Unit]
Description=Read Route Record 1.5

[Service]
ExecStart=/usr/bin/rrr -s -l /etc/rrr.conf.d/
KillSignal=SIGTERM
User=daemon
Group=daemon

[Install]
WantedBy=multi-user.target
