# misc/systemd/rrr.service
[Unit]
Description=Read Route Record 0.19

[Service]
ExecStart=/usr/bin/rrr -l /etc/rrr.conf.d/
KillSignal=SIGINT
User=rrr
Group=rrr

[Install]
WantedBy=multi-user.target
