# misc/docker/service/services.d/rrr-logd/finish
#!/bin/sh

echo "Service rrr-logd exited with code $1, restarting"

exit 0
