# misc/docker/service/services.d/rrr-logd/run
#!/bin/sh -e

cd /opt/app

rm -f $LOG_SOCKET

rrr_logd -s $LOG_SOCKET -p -n -u
